# docker/prometheus/node_exporter.service
[Unit]
Description=Node Exporter

[Service]
User=root
ExecStart=/node_exporter

[Install]
WantedBy=default.target
